# launchd/com.artofdrawers.opsdashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AoD Operations Dashboard — launchd schedule
  ============================================
  Replaces the Cowork scheduled task (which ran in a sandbox that couldn't
  write to /tmp, couldn't touch .git/, and didn't have xlrd installed).

  Install:
    mkdir -p ~/Library/LaunchAgents
    cp launchd/com.artofdrawers.opsdashboard.plist ~/Library/LaunchAgents/
    launchctl unload ~/Library/LaunchAgents/com.artofdrawers.opsdashboard.plist 2>/dev/null
    launchctl load   ~/Library/LaunchAgents/com.artofdrawers.opsdashboard.plist

  Inspect:
    launchctl list | grep opsdashboard         # confirm it's loaded
    tail -f ops-dashboard/logs/run.log          # watch live runs
    tail -f ops-dashboard/logs/run.err          # watch errors

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.artofdrawers.opsdashboard.plist
    rm ~/Library/LaunchAgents/com.artofdrawers.opsdashboard.plist

  Schedule: M-F at 9:00am, 12:00pm, 3:00pm local time.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.artofdrawers.opsdashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/artofdrawersllc/Documents/Claude/Projects/AoD_Cowork/ops-dashboard/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/artofdrawersllc/Documents/Claude/Projects/AoD_Cowork/ops-dashboard</string>

    <!--
      Weekday: 1=Mon, 2=Tue, 3=Wed, 4=Thu, 5=Fri  (0/7=Sun, 6=Sat)
      One <dict> per (weekday, hour). 5 days × 3 times = 15 entries.
    -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <!-- Append-only log files so you can see exactly what each run did. -->
    <key>StandardOutPath</key>
    <string>/Users/artofdrawersllc/Documents/Claude/Projects/AoD_Cowork/ops-dashboard/logs/run.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/artofdrawersllc/Documents/Claude/Projects/AoD_Cowork/ops-dashboard/logs/run.err</string>

    <!-- Don't run immediately on load — wait for the next scheduled slot. -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
